# supervisor.conf
[program:mysqld]
user=mysql
command=/usr/bin/mysqld_safe
redirect_stderr=true
autorestart=true

[program:apache]
command=/apache_foreground.sh
redirect_stderr=true
autorestart=true
